# da_le_tou/crontab
0 18 * * *  python3 /hlp/fund/da_le_tou/store_lucky_numbers.py >> /hlp/fund/da_le_tou/lucky_numbers.log 2>&1
