# da_le_tou/crontab
0 18 * * *  python3 /hlp/fund/store_lucky_numbers.py >> /hlp/fund/lucky_numbers.log 2>&1